# ansible/files/etc/systemd/system/rhythm.service
[Unit]
Description=TanStack Service (Rhythm Place)
After=network.target

[Service]
Type=simple
User=nginx
WorkingDirectory=/var/www/rhythm.place
ExecStart=/usr/bin/npm run start
Restart=on-failure

[Install]
# start the app automatically
WantedBy=multi-user.target