# files/rhythm-autoplay.service
[Unit]
Description=Auto-resume Spotify playback on Rhythm Place
After=rhythm-liquidsoap.service
Requires=rhythm-liquidsoap.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rhythm-autoplay.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target